# backend/crontab
SHELL=/bin/bash

*/1 * * * * echo "Writing data to redis..." > /proc/1/fd/1 2>/proc/1/fd/2
*/1 * * * * source $HOME/.profile; python /get_data.py >> /var/log/cron.log 2>&1

